# config/systemd/user/omarchy-chromium-shutdown.service
[Unit]
Description=Graceful Chromium-based browser shutdown
Documentation=https://wiki.archlinux.org/title/Systemd/User
DefaultDependencies=no
Before=wayland-session-shutdown.target
PartOf=wayland-session@.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true
TimeoutStopSec=15
ExecStop=%h/.local/bin/omarchy-chromium-shutdown

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=wayland-session@.target
