# config/systemd/user/omarchy-chromium-shutdown.service
[Unit]
Description=Graceful Chromium-based browser shutdown
DefaultDependencies=no

Before=wayland-wm@.service
Before=wayland-wm-env@.service
Before=wayland-session-shutdown.target
Conflicts=wayland-session-shutdown.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/omarchy-chromium-shutdown
TimeoutStartSec=20
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=wayland-session-shutdown.target
